# bin/tunnel-server.service
[Unit]
Description=Tunnel Server Service
After=network.target
 
[Service]
Type=simple
WorkingDirectory=/opt/tunnel
ExecStart=/opt/tunnel/tunnel server
KillMode=control-group
Restart=on-failure
LimitNOFILE=100000
 
[Install]
WantedBy=multi-user.target